# deploy/portfolio.service
[Unit]
Description=Portfolio Service
After=network.target mariadb.service
Requires=mariadb.service

[Service]
Type=simple
WorkingDirectory=/opt/portfolio
Environment=PYTHONPATH=/opt/portfolio
EnvironmentFile=/etc/portfolio.env
ExecStart=/opt/portfolio/venv/bin/python /opt/portfolio/service/start.py --host ${API_HOST} --port ${API_PORT}
Restart=on-failure
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
