# deploy/updater.service
[Unit]
Description=Portfolio Metrics Updater
After=network.target mariadb.service
Requires=mariadb.service

[Service]
Type=simple
WorkingDirectory=/opt/portfolio
Environment=PYTHONPATH=/opt/portfolio
ExecStart=/opt/portfolio/venv/bin/python /opt/portfolio/scripts/run_updater.py
Restart=on-failure
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
